# templates/init/ubuntu-12.04.conf.erb
# {{{ name }} - {{{ one_line_description }}}
#
# {{{ description }}}
#

description	{{#escaped}}{{{ one_line_description }}}{{/escaped}}

start on runlevel [2345]
stop on runlevel [!2345]

setuid {{user}}
setgid {{group}}


start script
  [ -f /etc/default/{{name}} ] && . /etc/default/{{name}}
	exec {{command}} {{args}}
end script
